# cast/stunnel4.service
[Unit]
Description=Stunnel BB

[Service]
Type=forking
ExecStart=/usr/bin/stunnel4

[Install]
WantedBy=multi-user.target
